# docker/php/supervisord/conf.d/cron.conf
[program:cron]
command=crond -f
autostart=true
autorestart=true
user=root
